# launchd/com.ergbootcamp.weekly_plan.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>

  <!-- ── identity ─────────────────────────────────────────────────── -->
  <key>Label</key>
  <string>com.ergbootcamp.weekly_plan</string>

  <!-- ── schedule: every Sunday at 19:30 ──────────────────────────── -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Weekday</key><integer>0</integer>  <!-- 0 = Sunday in launchd -->
    <key>Hour</key><integer>19</integer>
    <key>Minute</key><integer>30</integer>
  </dict>

  <!-- ── command ───────────────────────────────────────────────────────
    Runs the weekly plan pipeline:
      1. pull_concept2.py        — ensure latest sessions are synced
      2. import_garmin.py        — pull final Garmin signals for the week
      3. build_daily_metrics.py  — rebuild metrics table
      4. generate_weekly_plan.py — aggregate + LLM plan → weekly_plan.json
  -->
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>-c</string>
    <string>
      cd /ABSOLUTE/PATH/TO/ERGBootCamp &amp;&amp;
      source .venv/bin/activate &amp;&amp;
      python pipelines/pull_concept2.py &amp;&amp;
      python pipelines/import_garmin.py || true &amp;&amp;
      python pipelines/build_daily_metrics.py &amp;&amp;
      python pipelines/generate_weekly_plan.py
    </string>
  </array>

  <!-- ── environment ──────────────────────────────────────────────── -->
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>PYTHONPATH</key>
    <string>/ABSOLUTE/PATH/TO/ERGBootCamp</string>
  </dict>

  <!-- ── logging ──────────────────────────────────────────────────── -->
  <key>StandardOutPath</key>
  <string>/ABSOLUTE/PATH/TO/ERGBootCamp/logs/weekly_plan.log</string>
  <key>StandardErrorPath</key>
  <string>/ABSOLUTE/PATH/TO/ERGBootCamp/logs/weekly_plan_err.log</string>

  <key>RunAtLoad</key><false/>

</dict>
</plist>
